# scim.socket
[Unit]
Description=Input Service Framework Panel

[Socket]
ListenStream=/tmp/scim-socket-frontend
SmackLabelIPIn=*
SmackLabelIPOut=@
